# cookbooks/app_tomcat/templates/default/tomcat_init.erb
#!/bin/bash
# description: Tomcat Start Stop Restart
# processname: tomcat
# chkconfig: 234 20 80

TOMCAT_PROG="$(basename $0)"
NAME="$TOMCAT_PROG"
JAVA_HOME=$(readlink /usr/java/latest)
export JAVA_HOME
PATH=$JAVA_HOME/bin:$PATH
export PATH
CATALINA_HOME=/usr/share/apache-tomcat-<%=@version%>

if [ -r "$CATALINA_BASE/bin/setenv.sh" ]; then
  . "$CATALINA_BASE/bin/setenv.sh"
elif [ -r "$CATALINA_HOME/bin/setenv.sh" ]; then
  . "$CATALINA_HOME/bin/setenv.sh"
fi

start(){
  echo -n $"Starting $TOMCAT_PROG: "
  echo ""
  /bin/sh $CATALINA_HOME/bin/startup.sh
  RETVAL="$?"
  echo
  [ $RETVAL -eq 0 ] && touch /var/lock/subsys/$NAME
  return $RETVAL
}

stop(){
  echo -n $"Stopping $TOMCAT_PROG: "
  echo  ""
  
  /bin/sh $CATALINA_HOME/bin/shutdown.sh
  RETVAL="$?"

  if [ $RETVAL -eq 0 ] ; then
    rm -f /var/lock/subsys/$NAME
  else
    echo "Unable to stop $NAME"
  fi
  echo
  return $RETVAL
}

restart(){
  stop
  sleep 10
  start
}

version(){
  $CATALINA_HOME/bin/version.sh
}

case $1 in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    restart
    ;;

  version)
    version
    ;;
  *)
    echo $"Usage: $0 {start|stop|restart|version}"
    exit 1
esac
exit 0
